# trading-bot/systemd/trading-bot-web.service
[Unit]
Description=Oil Trading Bot Web Dashboard (venv)
After=network.target

[Service]
Type=simple
User=claw
WorkingDirectory=/home/claw/.openclaw/workspace/trading-bot
Environment="PATH=/home/claw/.openclaw/workspace/trading-bot/venv/bin"
ExecStart=/home/claw/.openclaw/workspace/trading-bot/venv/bin/python /home/claw/.openclaw/workspace/trading-bot/web/server.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
